# konemqsend.conf
description "My service"
author  "Some Dude <blah@foo.com>"

start on runlevel [234]
stop on runlevel [0156]

chdir /home/pi/
exec /home/pi/script.py
respawn
